# deploy/supervisord_recovery.conf
[program:recovery_service]
command=/path/to/venv/bin/python /path/to/Database-Track-2/core/recovery_service.py --interval 60
directory=/path/to/Database-Track-2
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/recovery_service.out
stderr_logfile=/var/log/recovery_service.err
# set user to run the process
user=REPLACE_WITH_USER
